# packaging/aur/whisper-agent-sandbox.service
[Unit]
Description=whisper-agent sandbox daemon (per-thread landlock-isolated MCP hosts)
Documentation=https://github.com/computer-whisperer/whisper-agent
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
# Runs as root by design. The daemon's job is to provision
# landlock-isolated mcp-host children that reach into host files on
# behalf of a remote whisper-agent server; confining the daemon under
# an unprivileged service user would subtract DAC reach before
# landlock ever narrows it, leaving the gate mostly-closed regardless
# of what the pod's `allowed_paths` lists. Landlock is the gate.
# (System-level hardening below still applies to the daemon process.)

# Tunables (LISTEN address etc.). Optional; the unit ships with a
# sample at /etc/whisper-agent/sandbox.env populated with the
# defaults.
EnvironmentFile=-/etc/whisper-agent/sandbox.env
Environment=WHISPER_AGENT_SANDBOX_LISTEN=[::]:9820

# systemd creates /var/lib/whisper-agent before ExecStartPre fires, so
# the helper can drop the control token there without a separate
# tmpfiles.d entry.
StateDirectory=whisper-agent
StateDirectoryMode=0750

# Generate the control-plane bearer token on first start. Idempotent:
# the helper no-ops when the file is already populated. Rotate by
# deleting the file and restarting; remote whisper-agent servers
# holding the old token will then fail auth and need re-provisioning.
ExecStartPre=/usr/lib/whisper-agent/ensure-control-token /var/lib/whisper-agent/sandbox-control-token

ExecStart=/usr/bin/whisper-agent-sandbox \
    --listen ${WHISPER_AGENT_SANDBOX_LISTEN} \
    --mcp-host-bin /usr/bin/whisper-agent-mcp-host \
    --control-token-file /var/lib/whisper-agent/sandbox-control-token

Restart=on-failure
RestartSec=5

# --- Hardening ---------------------------------------------------------
# The daemon runs with its default (full) root capability set. It
# needs CAP_DAC_READ_SEARCH + CAP_DAC_OVERRIDE to open user-owned
# paths at restrictive modes like 0700 — those fds become the input
# to landlock's `path_beneath_rules`. Stripping DAC caps here would
# leave landlock silently missing rules for any workspace under a
# user's home. The spawned mcp-host child then operates within the
# landlock restrictions regardless of the caps it inherited.
NoNewPrivileges=yes

# Filesystem isolation.
#
# Deliberately NOT setting PrivateTmp=yes: pods configure their
# workspace root as a host path (commonly under /tmp for scratch
# work), and the child MCP host inherits this service's mount
# namespace. A private /tmp would silently shadow those paths — the
# child would see an empty tmpfs and fail to start. Per-thread
# isolation is landlock's job, not systemd's.
#
# ProtectHome is NOT set here either: pods routinely declare
# workspaces under `/home/<user>/` (e.g. a data dir, a project
# checkout), and landlock — applied per-child with the pod's
# `allowed_paths` — is what actually scopes the child's reach. Having
# systemd blanket-deny /home at the daemon level would prevent the
# very resolution the child needs to perform.
ProtectSystem=strict
ReadWritePaths=/var/lib/whisper-agent

# Kernel surface: nothing the daemon needs.
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Process surface
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
