# packaging/aur/whisper-agent-host.install
post_install() {
    cat <<MSG
==> whisper-agent-host installed.

  This package ships the host-side daemons:
    /usr/bin/whisper-agent-sandbox    landlock-isolated MCP host provisioner
    /usr/bin/whisper-agent-mcp-host   the per-thread MCP server (spawned by sandbox)

  The daemon runs as root. Landlock rules — applied per-child from
  each pod's \`allowed_paths\` — are the filesystem gate. Running as a
  dedicated service user would subtract DAC reach before landlock
  ever narrowed it, leaving workspaces like /home/<you>/projects
  unreachable regardless of what the pod config requested.

  The control-plane bearer token is auto-generated on first start at:
    /var/lib/whisper-agent/sandbox-control-token   (mode 0600)

  Edit /etc/whisper-agent/sandbox.env to change the listen address
  (default: [::]:9820, dual-stack — leaves 9810 free for the in-tree
  dev harness at scripts/dev.sh).

  Enable and start the sandbox daemon:
    systemctl enable --now whisper-agent-sandbox.service

  After first start, copy the control token to your whisper-agent
  server (e.g. into the kubernetes secret it mounts as
  --host-env-provider-token name=path) so the server can authenticate.
MSG
}

post_upgrade() {
    # Pre-0.3.0-2 the daemon ran under a dedicated 'whisper-agent'
    # system user. That model is gone — landlock alone is the gate
    # now. Old installs may still have the user + the root-owned
    # state dir; clean up best-effort so a restart picks up the new
    # user cleanly. Failure here is non-fatal (the admin may have
    # repurposed the user).
    if getent passwd whisper-agent >/dev/null; then
        echo "==> removing legacy 'whisper-agent' system user (no longer used)"
        userdel whisper-agent 2>/dev/null || true
    fi
    if getent group whisper-agent >/dev/null; then
        groupdel whisper-agent 2>/dev/null || true
    fi
    # Re-chown the state dir so the root-owned daemon can write it.
    if [[ -d /var/lib/whisper-agent ]]; then
        chown -R root:root /var/lib/whisper-agent 2>/dev/null || true
    fi

    if systemctl is-active --quiet whisper-agent-sandbox.service 2>/dev/null; then
        echo "==> restart whisper-agent-sandbox.service to pick up new binaries"
    fi
}

pre_remove() {
    if systemctl is-active --quiet whisper-agent-sandbox.service 2>/dev/null; then
        systemctl stop whisper-agent-sandbox.service
    fi
    if systemctl is-enabled --quiet whisper-agent-sandbox.service 2>/dev/null; then
        systemctl disable whisper-agent-sandbox.service
    fi
}

post_remove() {
    cat <<MSG
==> whisper-agent-host removed.

  The runtime data directory /var/lib/whisper-agent/ was preserved
  (sandbox-control-token, pods/, audit.jsonl). Remove it manually
  if no longer needed:

    rm -rf /var/lib/whisper-agent
MSG
}
