# packaging/aur/whisper-agent-host-daemon.service
[Unit]
Description=whisper-agent host-env daemon (dials a whisper-agent server, provisions per-thread landlock-isolated MCP hosts)
Documentation=https://github.com/computer-whisperer/whisper-agent
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
# Runs as root by design. The daemon's job is to provision
# landlock-isolated mcp-host children that reach into host files on
# behalf of a remote whisper-agent server; confining the daemon under
# an unprivileged service user would subtract DAC reach before
# landlock ever narrows it, leaving the gate mostly-closed regardless
# of what the pod's `allowed_paths` lists. Landlock is the gate.
# (System-level hardening below still applies to the daemon process.)

# Required tunables — admin must set these in /etc/whisper-agent/host-daemon.env
# before enabling the unit:
#   WHISPER_AGENT_SERVER_URL  — wss://your-server.tld/v1/host_env_link
#   WHISPER_AGENT_TOKEN_FILE  — file containing the bearer token; the
#                               same value must appear in [[auth.daemons]]
#                               on the central server.
EnvironmentFile=/etc/whisper-agent/host-daemon.env

# State directory still useful for caching things the daemon may want
# across restarts (and as a per-host workspace root for the startup
# probe). systemd creates /var/lib/whisper-agent before ExecStart fires.
StateDirectory=whisper-agent
StateDirectoryMode=0750

ExecStart=/usr/bin/whisper-agent-host-daemon \
    --server-url ${WHISPER_AGENT_SERVER_URL} \
    --token-file ${WHISPER_AGENT_TOKEN_FILE} \
    --mcp-host-bin /usr/bin/whisper-agent-mcp-host \
    --probe-workspace /var/lib/whisper-agent

Restart=on-failure
RestartSec=5

# --- Hardening ---------------------------------------------------------
# The daemon runs with its default (full) root capability set. It
# needs CAP_DAC_READ_SEARCH + CAP_DAC_OVERRIDE to open user-owned
# paths at restrictive modes like 0700 — those fds become the input
# to landlock's `path_beneath_rules`. Stripping DAC caps here would
# leave landlock silently missing rules for any workspace under a
# user's home. The spawned mcp-host child then operates within the
# landlock restrictions regardless of the caps it inherited.
NoNewPrivileges=yes

# Filesystem isolation.
#
# Deliberately NOT setting PrivateTmp=yes: pods configure their
# workspace root as a host path (commonly under /tmp for scratch
# work), and the child MCP host inherits this service's mount
# namespace. A private /tmp would silently shadow those paths — the
# child would see an empty tmpfs and fail to start. Per-thread
# isolation is landlock's job, not systemd's.
#
# ProtectHome is NOT set here either: pods routinely declare
# workspaces under `/home/<user>/` (e.g. a data dir, a project
# checkout), and landlock — applied per-child with the pod's
# `allowed_paths` — is what actually scopes the child's reach. Having
# systemd blanket-deny /home at the daemon level would prevent the
# very resolution the child needs to perform.
ProtectSystem=strict
ReadWritePaths=/var/lib/whisper-agent

# Kernel surface: nothing the daemon needs.
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Process surface
LockPersonality=yes
RestrictRealtime=yes
# RestrictSUIDSGID is deliberately off: the daemon's worker spawn
# drops uid to the pod-binding's `runas` user via setgroups + setgid
# + setuid in the forked child's pre_exec, before landlock seals it.
# RestrictSUIDSGID=yes would kill the worker with SIGSYS on the very
# first of those calls. The daemon itself stays root; the carve-out
# only enables the legitimate drop-uid path on each spawn.
RestrictSUIDSGID=no
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
# Strip @privileged and @resources, then add back exactly the three
# privileged calls the worker spawn's drop-uid chain needs
# (setgroups / setgid / setuid). Everything else in @privileged
# (mount, swapon, ptrace targets, etc.) stays denied.
SystemCallFilter=~@privileged @resources
SystemCallFilter=setgroups setgid setuid

[Install]
WantedBy=multi-user.target
